# utility/roles/aptly/templates/aptly-ws.conf.j2
# aptly-ws.conf
start on filesystem
exec aptly serve -listen=0.0.0.0:8082