# utility/roles/consul/templates/consul-ui.conf.j2
# consul-ui.conf
start on filesystem
exec /usr/lib/consul agent -data-dir /tmp/consul-ui -client 192.168.10.11 -ui-dir /opt/consul-ui -join 192.168.10.10 -bind 192.168.10.11 -node=consul-ui
