# utility/roles/consul/templates/consul.conf.j2
# consul.conf
start on filesystem
exec /usr/lib/consul agent -server -bootstrap -bind=192.168.10.10 -data-dir /tmp/consul